# tocsin.service
[Unit]
Description=Plays a chime on the quarter hour

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tocsin
